# deploy/magnetic-platform.service
[Unit]
Description=Magnetic Platform Server (multi-tenant V8 hosting)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/magnetic-v8-server --platform --port 3003 --data-dir /var/lib/magnetic/apps --park-idle 300
Restart=always
RestartSec=2
LimitNOFILE=65535

Environment=RUST_LOG=info

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/magnetic

[Install]
WantedBy=multi-user.target
